# apps/opensync-ubuntu/packaging/stage/opensync-ubuntu_0.1.0_amd64/lib/systemd/user/opensync-ubuntu.service
[Unit]
Description=OpenSync Ubuntu (vault sync)
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node /usr/lib/opensync-ubuntu/dist/cli.js run
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
